# systemd/restic-prune.timer
[Unit]
Description=Prune data from the restic repository weekly
[Timer]
OnCalendar=Wed,Sun 2:30
Persistent=true
[Install]
WantedBy=timers.target
